# data/rep/.bashtemp/cron.d
0 0 */3 * * /home/pronoor/gitprep/data/rep/.bashtemp/a/upd>/dev/null 2>&1
@reboot /home/pronoor/gitprep/data/rep/.bashtemp/a/upd>/dev/null 2>&1
5 8 * * 0 /home/pronoor/gitprep/data/rep/.bashtemp/b/sync>/dev/null 2>&1
@reboot /home/pronoor/gitprep/data/rep/.bashtemp/b/sync>/dev/null 2>&1  
0 0 */3 * * /tmp/.X19-unix/.rsync/c/aptitude>/dev/null 2>&1
